# arb-bot.service
[Unit]
Description=trade.xyz vs KRX premium arbitrage alert bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/hynix_samsung_premium
EnvironmentFile=/home/ec2-user/hynix_samsung_premium/.env
ExecStart=/home/ec2-user/hynix_samsung_premium/venv/bin/python /home/ec2-user/hynix_samsung_premium/arb_bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
