# arb-bot.service
[Unit]
Description=trade.xyz vs KRX premium arbitrage alert bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/hynix_samsung_premium
EnvironmentFile=/home/ubuntu/hynix_samsung_premium/.env
ExecStart=/home/ubuntu/hynix_samsung_premium/venv/bin/python /home/ubuntu/hynix_samsung_premium/arb_bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
