# arb-bot.service
[Unit]
Description=trade.xyz vs KRX premium arbitrage alert bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/sk_samsung_funding
EnvironmentFile=/home/ubuntu/sk_samsung_funding/.env
ExecStart=/home/ubuntu/sk_samsung_funding/venv/bin/python /home/ubuntu/sk_samsung_funding/arb_bot.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
